# root/etc/service/cron/run
#!/bin/sh -e
exec 2>&1

rm -f /var/run/crond.pid

exec cron -f
